# roles/cart/templates/cart.service.j2
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis.nanda.cfd
Environment=CATALOGUE_HOST=catalogue.nanda.cfd
Environment=CATALOGUE_PORT=8080
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target